# scripts/orchestration/systemd/odds_worker.service
[Unit]
Description=Football Predictor Odds Worker
After=network.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/football_predictor
ExecStart=/opt/football_predictor/.venv/bin/python scripts/odds_worker.py
StandardOutput=append:/opt/football_predictor/logs/odds_worker.log
StandardError=append:/opt/football_predictor/logs/odds_worker.log
